# cmd/shuttle-tunnel.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>{{ .Label }}</string>

  <key>ProgramArguments</key>
  <array>
    {{- if .Multiplex }}
    {{- /* 2FA hosts (nibi's Duo): the tunnel can only ride a
           human-authenticated ControlMaster socket — a headless job can never
           answer keyboard-interactive, and plain `ssh -S <sock>` silently
           falls back to a FRESH connection when the socket is absent (Duo
           denial spam, then autossh exits 0 and SuccessfulExit=false parks
           the job dead). So: supervise with a shell loop that gates every
           attempt on `ssh -O check` — socket dead → poll every 30s, quiet;
           socket alive → run the forward (BatchMode=yes so a race can never
           degrade into interactive auth). Master dies → ssh exits → back to
           polling. The next approved human login revives the tunnel on its
           own, within 30s. */}}
    <string>/bin/sh</string>
    <string>-c</string>
    <string>while :; do if ssh -O check {{ .SSHHost }} 2>/dev/null; then ssh -N -S {{ .Home }}/.ssh/ctl/%C -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -o ControlMaster=no -o ExitOnForwardFailure=yes -o BatchMode=yes -L {{ .LocalPort }}:localhost:4000 {{ .SSHHost }}; fi; sleep 30; done</string>
    {{- else }}
    <string>{{ .AutoSSHPath }}</string>
    {{- if not .HoldCommand }}
    <string>-N</string>
    {{- end }}
    <string>-M</string>
    <string>0</string>
    <string>-S</string>
    <string>none</string>
    <string>-o</string>
    <string>ServerAliveInterval=30</string>
    <string>-o</string>
    <string>ServerAliveCountMax=3</string>
    <string>-o</string>
    <string>ControlMaster=no</string>
    <string>-o</string>
    <string>ExitOnForwardFailure=yes</string>
    {{- if .SSHAuthSock }}
    <string>-o</string>
    <string>IdentityAgent={{ .SSHAuthSock }}</string>
    {{- end }}
    <string>-L</string>
    <string>{{ .LocalPort }}:localhost:4000</string>
    <string>{{ .SSHHost }}</string>
    {{- if .HoldCommand }}
    <string>{{ .HoldCommand }}</string>
    {{- end }}
    {{- end }}
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>AUTOSSH_GATETIME</key>
    <string>0</string>
    {{- if .SSHAuthSock }}
    <key>SSH_AUTH_SOCK</key>
    <string>{{ .SSHAuthSock }}</string>
    {{- end }}
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>NetworkState</key>
    <true/>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>{{ .LogPath }}</string>
  <key>StandardErrorPath</key>
  <string>{{ .LogPath }}</string>
  <key>WorkingDirectory</key>
  <string>{{ .Home }}</string>
</dict>
</plist>
